# build/setup/usersvc.service
[Unit]
Description=Usersvc
Wants=network-online.target
After=network-online.target

[Service]
User=gradtage
Group=gradtage
Type=simple
ExecStartPre=/bin/sleep 1
ExecStart=/usr/local/bin/usersvc \
--config.file /etc/gradtage/usersvc.yml
Restart=always
RestartSec=3


[Install]
WantedBy=multi-user.target